# nhl-fetcher.service
[Unit]
Description=NHL Game Results Fetcher
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/workspace
Environment=DJANGO_SETTINGS_MODULE=bet_bot.settings.production
ExecStart=/usr/bin/python3 /workspace/manage.py fetch_nhl_scores
StandardOutput=journal
StandardError=journal